# conf/supervisor.ini
[fcgi-program:pipages]
socket=tcp://localhost:10234
command=/opt/pipages/venv/bin/daphne -u /opt/data/pipages/pipages%(process_num)d.sock --fd 0 --access-log - --proxy-headers pipages.asgi:application
directory=/opt/pipages
numprocs=4
process_name=pipages%(process_num)d
user=root
autostart=true
autorestart=false
stdout_logfile=/opt/pipages/logs/pipages.log
redirect_stderr=true

[program:celery_beat]
command=/opt/pipages/venv/bin/celery -A pipages beat -l info
directory=/opt/pipages
user=root
autostart=true
autorestart=true
stdout_logfile=/opt/pipages/logs/celery_beat.log

[program:celery_worker]
command=/opt/pipages/venv/bin/celery -A pipages worker -l info
directory=/opt/pipages
user=root
autostart=true
autorestart=true
stdout_logfile=/opt/pipages/logs/celery_worker.log

